# s6-rc.d/cron/run
#!/command/execlineb -P

with-contenv
cron -f -L /proc/1/fd/1
